# crontab
*/20 * * * * bundle exec rails tournaments:update
